# .crontab
1 * * * * ~/git-update.sh >> git-update.log