# dot/.crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.QLD0da installed on Thu Jul 28 15:31:16 2016)
# (Cronie version 4.2)
0 16 * * * ~/github/os-script/sh/backup.sh
0 16 * * * python ~/github/os-script/py/classify.py


